# logging/server/lamassu-server-prof.conf
[program:lamassu-server]
command=node --prof /usr/local/bin/lamassu-server
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/lamassu-server.err.log
stdout_logfile=/var/log/supervisor/lamassu-server.out.log
environment=HOME="/root"
