# wallets/repair/orig/zcash.conf
[program:zcash]
command=nice /usr/local/bin/zcashd -datadir=/mnt/blockchains/zcash
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/zcash.err.log
stdout_logfile=/var/log/supervisor/zcash.out.log
stderr_logfile_backups=2
stdout_logfile_backups=2
environment=HOME="/root"
